# root/usr/lib/systemd/system/nethserver-config-network.service
[Unit]
Description=Reconfigure newtork interfaces
Documentation=https://github.com/NethServer/nethserver-base
ConditionPathExists=!/var/spool/first-boot
ConditionFileNotEmpty=/var/lib/nethserver/db/networks
After=network-online.target
RefuseManualStart=true
RefuseManualStop=true

[Service]
Type=oneshot
Environment=SKIPCARRIER=1
ExecStart=/sbin/e-smith/nethserver-config-network

[Install]
WantedBy=multi-user.target
