# root/etc/cron.daily/conf-mod_ssl
#! /bin/sh

exec /sbin/e-smith/expand-template /home/e-smith/ssl.pem/pem
