# root/etc/cron.daily/nethserver-check-builtin-certificate
#!/bin/bash

#
# Copyright (C) 2013 Nethesis S.r.l.
# http://www.nethesis.it - support@nethesis.it
#
# This script is part of NethServer.
#
# NethServer is free software: you can redistribute it and/or modify
# it under the terms of the GNU General Public License as published by
# the Free Software Foundation, either version 3 of the License,
# or any later version.
#
# NethServer is distributed in the hope that it will be useful,
# but WITHOUT ANY WARRANTY; without even the implied warranty of
# MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the
# GNU General Public License for more details.
#
# You should have received a copy of the GNU General Public License
# along with NethServer.  If not, see <http://www.gnu.org/licenses/>.
#

program=$(basename $0)
pki_info=/usr/libexec/nethserver/pki-info

# 10 days advance alert
alert=10

# Have custom certificate?
cert=/etc/pki/tls/certs/NSRV.crt

# Certificate expire timestamp
exp=$(${pki_info} -e ${cert})

# Renew, if cert is stale and not custom
must_renew() {
    perl -MDate::Parse -e '
$oneday = str2time($ARGV[0]) < (time() + (24 * 3600));
exit($oneday ? 0 : 1);
' -- "${exp}"
}

#
# Send mail to root user
#

if must_renew; then
    {
	/etc/e-smith/events/actions/nethserver-generate-certificate
	echo "$(hostname) certificate has been automatically generated"
	echo ""
	${pki_info} -s ${cert} 2>&1
    } | mail -s "$(hostname) certificate has been generated" root@localhost
    echo "[WARNING] $(hostname) certificate has been generated" | logger -t "${program}" -p user.warn
fi
