# root/etc/cron.daily/nethserver-letsencrypt-certs
#!/bin/bash

/usr/libexec/nethserver/letsencrypt-certs
