# root/etc/udev/rules.d/95-NetworksDB.rules
#
# Update NetworksDB 
#
SUBSYSTEM=="net", ACTION=="*", KERNEL=="*", RUN+="/sbin/e-smith/nethserver-config-network --update-db-only"
