# root/etc/udev/rules.d/95-NetworksDB.rules
#
# Update NetworksDB 
#
SUBSYSTEM=="net", ACTION=="*", KERNEL=="*", RUN+="/usr/libexec/nethserver/udev-update-network %p %k"
