# root/etc/udev/rules.d/95-NetworksDB.rules
#
# Update NetworksDB 
#
SUBSYSTEM=="net", ACTION=="*", KERNEL=="*", RUN+="/usr/bin/systemctl start nethserver-config-network.service"
