# root/etc/cron.weekly/refresh-ddns
#! /bin/bash

exec /etc/e-smith/events/actions/update-dns
