# root/etc/cron.d/warnpassexpire
# Send email reminders about change the password
42 0 * * *   root /sbin/e-smith/warnpassexpire
